# postfix.cron
#!/bin/sh
#
# Check for problems with Postfix

/usr/sbin/postfix check